# wordle/wordle.service
# /etc/systemd/system/wordle.service

[Unit]
Description=Wordle
After=multi-user.target

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/python3 /home/user/wordle.py

[Install]
WantedBy=multi-user.target
